# systemd/.config/systemd/user/agents.target
[Unit]
Description=Collection of all agents that need to run for me.