# systemd/.config/systemd/user/mate-utils.target
[Unit]
Description=Collects all mate daemons that are needed when running exwm. Is activated in .xinitrc.
